# build/COPY_ROOT_0/etc/services.d/dbus/run
#!/bin/bash

set -e # Exit immediately if a command exits with a non-zero status.
set -u # Treat unset variables as an error.

source /opt/ai-dock/etc/environment.sh

sudo mkdir -p /var/run/dbus
sudo mkdir -p /run/dbus  # Ensure /run/dbus directory exists

# Generate machine ID if it doesn't exist
if [ ! -f /var/lib/dbus/machine-id ]; then
    sudo dbus-uuidgen > /var/lib/dbus/machine-id
fi

# Clean up any old sockets
sudo rm -f /var/run/dbus/system_bus_socket
sudo rm -f /run/dbus/system_bus_socket
sudo rm -f /run/dbus/pid

# Start system-level D-Bus daemon
echo "Starting system D-Bus daemon..."
sudo dbus-daemon --system --nosyslog --nofork --print-address &

# Wait for D-Bus socket creation
echo "Waiting for D-Bus socket creation..."
MAX_WAIT=10
count=0
while [ ! -S /run/dbus/system_bus_socket ] && [ $count -lt $MAX_WAIT ]; do
    sleep 1
    count=$((count+1))
    echo "Waiting for D-Bus socket... ($count/$MAX_WAIT)"
done

if [ ! -S /run/dbus/system_bus_socket ]; then
    echo "Error: D-Bus socket was not created"
    exit 1
fi

# Set DBUS_SYSTEM_BUS_ADDRESS environment variable
export DBUS_SYSTEM_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
echo "D-Bus system address: $DBUS_SYSTEM_BUS_ADDRESS"

# Wait for rtkit to start
sleep 2

# Test rtkitctl
rtkitctl --start

exec dbus-daemon --session --nosyslog --address=$DBUS_SESSION_BUS_ADDRESS --nofork

# vim:ft=sh:ts=4:sw=4:et:sts=4