# build/COPY_ROOT_0/etc/services.d/turbovnc/run
#!/bin/bash

set -e # Exit immediately if a command exits with a non-zero status.
set -u # Treat unset variables as an error.

sudo rm -rf /tmp/.X* ~/.cache

source /opt/ai-dock/etc/environment.sh

sudo mkdir -p /tmp/.X11-unix
sudo chmod 777 -R /tmp/.X11-unix

mkdir -p ~/.vnc/ ~/.dosbox
echo $USER_PASSWORD | /opt/TurboVNC/bin/vncpasswd -f | tee /root/.vnc/passwd > /dev/null
chmod 0600 ~/.vnc/passwd

exec /opt/TurboVNC/bin/vncserver "$@"